# scripts/crontab
*/5 *    * * * bash ${TT_HOME}/scripts/ps_checker.sh
